# examples/slapd-backup.cron.d
# slapd git backup
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
0 5 * * * root /usr/local/sbin/backup-slapd-git.sh
